# deploy/today_tactics.service
[Unit]
Description=Today Tactics — K리그 전술 분석
After=network.target

[Service]
User=rocky
Group=rocky
WorkingDirectory=/opt/today_tactics
Environment="PATH=/opt/today_tactics/venv/bin"
Environment="FLASK_DEBUG=0"

# CPU 우선순위 — 같은 서버의 다른 프로세스(bunting 봇) 폭주 시
# today_tactics가 사용자 응답 우선으로 CPU 받도록 비례 가중치 + nice 조정
# cgroup v2의 CPUWeight 기본값 100, 300은 약 3배 가중. Nice -5는 priority 약간 상향.
CPUWeight=300
Nice=-5
ExecStart=/opt/today_tactics/venv/bin/gunicorn \
    --workers 3 \
    --worker-class gthread \
    --threads 2 \
    --bind 127.0.0.1:5000 \
    --timeout 120 \
    --keep-alive 5 \
    --access-logfile /var/log/today_tactics/access.log \
    --error-logfile /var/log/today_tactics/error.log \
    main:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
